# ansible/roles/gnome/files/gnome-auto-dconf-dump.timer
[Unit]
Description=Run automatic GNOME settings dump

[Timer]
RandomizedDelaySec=120
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
