# modules/systemd/include/.config/systemd/user/playerctld.service
[Unit]
Description=Playerctl daemon
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/playerctld
Restart=on-failure
