# modules/systemd/include/.config/systemd/user/suspend.target
[Unit]
Description=User level suspend target (after returning from sleep)
StopWhenUnneeded=yes
Wants=screensaver-unlock.service
