# modules/twm-x11/include/.config/systemd/user/hibernate.target
[Unit]
Description=User level hibernate target (after returning from hibernation)
StopWhenUnneeded=yes
Wants=screensaver-unlock.service
