# modules/twm-x11/include/.config/systemd/user/screensaver-unlock.service
[Unit]
Description=Unlock user screen after sleeping
After=suspend.target

[Service]
Type=oneshot
ExecStart=%h/.bin/screen-locker unlock
